# deploy/supervisor/crm-queue.conf
; Supervisor config for the CRM outbound queue worker.
; Install: copy to /etc/supervisor/conf.d/crm-queue.conf then:
;   sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start crm-queue:*
;
; Why this exists: outbound message delivery (SendChannelMessageJob) and other
; work run on the redis queue (QUEUE_CONNECTION=redis on prod). Without a running
; worker, agent replies are created as QUEUED but never sent to the provider,
; i.e. "gửi tin không hoạt động". The deploy pipeline runs `queue:restart`,
; which only signals an already-running worker to reload — this file is what
; keeps it running.

[program:crm-queue]
process_name=%(program_name)s_%(process_num)02d
command=php8.4 /var/www/crm/artisan queue:work redis --sleep=1 --tries=3 --max-time=3600 --backoff=5
directory=/var/www/crm
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/crm/storage/logs/queue-worker.log
stopwaitsecs=3600
